# supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
loglevel=debug

[program:sshd]
command=/usr/sbin/sshd -D

[program:postgresql]
command=/sbin/entrypoint.sh
